# ytfp.service
[Unit]
Description=Youtube Lazy Watch Later ulauncher extension

[Service]
Type=oneshot
